# minictinit.service
[Unit]
Description=Mini Container Init Script

[Service]
ExecStart=/run/minictinit.sh start

[Install]
WantedBy=multi-user.target
